# data/keyledd.service
[Unit]
Description=keyledd daemon

[Service]
Type=simple
ExecStart=/usr/bin/keyledd

[Install]
WantedBy=multi-user.target
